# INSE6130/ids.conf
[program:ids_script]
command=python3 /ids_app/ids.py
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ids_script.out.log
stderr_logfile=/var/log/supervisor/ids_script.err.log
