# INSE6130/ids.conf
[program:ids_script]
command=python3 /ids_app/ids.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
